# config/supervisor.conf
[supervisord]
nodaemon=true
logfile=/app/log/supervisord.log

[program:backend-app]
command = java -jar /app/data/app.jar
autostart = true
autorestart = true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0

[program:nginx]
command = /usr/local/nginx/sbin/nginx -g "daemon off;"
autostart = true
autorestart = true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
